# v2ray-sub.service
# فایل سرویس systemd برای اجرای خودکار Flask
# مسیر: /etc/systemd/system/v2ray-sub.service

[Unit]
Description=V2Ray Subscription Management System
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/v2ray-sub
Environment="PATH=/home/v2ray-sub/venv/bin"
ExecStart=/home/v2ray-sub/venv/bin/gunicorn --workers 2 --bind 127.0.0.1:5000 app:app

# امنیت
PrivateTmp=true
NoNewPrivileges=true

# ری‌استارت خودکار در صورت خطا
Restart=always
RestartSec=10

# لاگ‌ها
StandardOutput=journal
StandardError=journal
SyslogIdentifier=v2ray-sub

[Install]
WantedBy=multi-user.target
